# alphatrader/systemd/binance-collector.service
[Unit]
Description=AlphaTrader Binance Data Collector
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=goku
Group=alphatrader
WorkingDirectory=/home/shared/alphatrader/scripts
ExecStart=/usr/bin/python3 binance_collector.py --daemon
Restart=always
RestartSec=30
StandardOutput=append:/var/log/alphatrader/binance-collector.log
StandardError=append:/var/log/alphatrader/binance-collector.log

# Environment
Environment="ZENKAI_DB_HOST=localhost"
Environment="ZENKAI_DB_NAME=zenkai_data"
Environment="ZENKAI_DB_USER=zenkai"
Environment="ZENKAI_DB_PASSWORD=zenkai_secure_2024"

[Install]
WantedBy=multi-user.target
